# detonation/deploy/detonation.slice
[Unit]
Description=Detonation Sandbox Resource Slice

[Slice]
CPUQuota=200%
MemoryMax=2G
MemoryHigh=1800M
TasksMax=512
